# scripts/ultimate-ascent.conf
start on runlevel [2345]
stop on runlevel [016]

exec /home/kinect/code/bin/UltimateAscent
pre-stop /home/kinect/code/bin/UltimateAscent
